# etc/cron.hourly/update_repo
#!/bin/dash

(cd /instance-config/ServoBot;git pull) 2> /dev/null > /dev/null
